# mushr_utils/startup/rc.local
#!/bin/bash -e

# TODO Add what the rc local file should be.
